# mystart.service
[Unit]
Description=My custom startup script

[Service]
ExecStart=/home/lufiaq/btproxy.sh start

[Install]
WantedBy=multi-user.target
